# deploy/systemd/brain.service
[Unit]
Description=intelligente-home brain (FastAPI)
After=network-online.target vllm.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/intelligente-home/mate-brain
# Prod override'ları repo dışı dosyadan (commit'li dev .env'i ezer: env-var > .env)
EnvironmentFile=/etc/intelligente-home/brain.env
ExecStart=/opt/intelligente-home/mate-brain/.venv/bin/python -m brain.main
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
